# .docker/crontab/today
* * * * * root echo "today is $(date)" >> /var/log/today.log

# end of crontab